# json-server.conf
[supervisord]
nodaemon=true
logfile=/appdata/log/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D

[program:json-server1]
command=json-server  /appdata/www/data.json

