# supervisord.conf
[supervisord]
nodaemon=true


[program:sshd]
command=/usr/sbin/sshd -D

[program:json-server1]
command=json-server  /appdata/www/data.json

